# site-cookbooks/genigames/templates/default/couchdb.conf.erb
#!upstart

description "genigames.dev.concord.org"
author      "rpk@pobox.com"

start on (local-filesystems and net-device-up IFACE=eth0)
stop  on shutdown

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

script
  exec sudo -u couchdb /home/deploy/build-couchdb/build/bin/couchdb
end script
