# site-cookbooks/genigames/templates/default/genigames.dev.concord.org.conf.erb
#!upstart

description "genigames.dev.concord.org"
author      "rpk@pobox.com"

start on (local-filesystems and net-device-up IFACE=eth0)
stop  on shutdown

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

script
  exec sudo -u www-data sh -c "NODE_PORT=80 NODE_ENV=production authbind --deep /usr/local/bin/coffee /var/www/app.coffee >> /var/www/log/app.log 2>> /var/www/log/error.log"
end script
